# scripts/yt-session/systemd/yt-keepwarm.service
[Unit]
Description=YouTube cookie-session keep-warm (reload the YT tab)
After=yt-chromium.service

[Service]
Type=oneshot
ExecStart=/opt/discord-soundboard/scripts/yt-session/bin/keepwarm.sh
